# home/.config/profile.d/30-path
# Include executables installed by cabal-install in PATH
if [ -d "$HOME/.cabal/bin" ] ; then
    PATH="$HOME/.cabal/bin:$PATH"
fi

# Include executables installed by stack in PATH
if [ -d "$HOME/.local/bin" ]; then
    PATH="$HOME/.local/bin:$PATH"
fi

# Include executables installed by cabin in PATH - takes precedence over
# cabal-install and stack
if [ -d "$HOME/.cabin/bin" ] ; then
    PATH="$HOME/.cabin/bin:$PATH"
fi

# set PATH so it includes user's private bin if it exists - has highest
# precedence
if [ -d "$HOME/bin" ] ; then
    PATH="$HOME/bin:$PATH"
fi

# Puts executables installed with current Javascript project in PATH
PATH="node_modules/.bin:$PATH"
