# home/.config/profile.d/30-path
# Include executables installed by cabal-install in PATH
PATH=$(prepend-to-path "$HOME/.cabal/bin")

# Include executables installed by stack in PATH
PATH=$(prepend-to-path "$HOME/.local/bin")

# Include executables installed by cabin in PATH - takes precedence over
# cabal-install and stack
PATH=$(prepend-to-path "$HOME/.cabin/bin")

# set PATH so it includes user's private bin if it exists - has highest
# precedence
PATH=$(prepend-to-path "$HOME/bin")

# Puts executables installed with current Javascript project in PATH
PATH=$(prepend-to-path "node_modules/.bin")

export PATH
