# home/.config/profile.d/30-path
# We can't use `prepend-to-path` here because it is in `$HOME/bin`
PATH="$HOME/bin:$PATH"

PATH=$(prepend-to-path "/usr/games")

# Include executables installed by cargo in PATH
PATH=$(prepend-to-path "$HOME/.cargo/bin")

# Include executables installed by cabal-install in PATH
PATH=$(prepend-to-path "$HOME/.cabal/bin")

# Include executables installed by stack in PATH
PATH=$(prepend-to-path "$HOME/.local/bin")

# Include executables installed by cabin in PATH - takes precedence over
# cabal-install and stack
PATH=$(prepend-to-path "$HOME/.cabin/bin")

# Include executables installed by yarn in PATH
PATH=$(prepend-to-path "$HOME/.config/yarn/global/node_modules/.bin")

# When working with Lerna projects, include top-level executables in PATH
PATH=$(prepend-to-path "../../node_modules/.bin")

# Puts executables installed with current Javascript project in PATH
PATH=$(prepend-to-path "node_modules/.bin")

# Use skim executables installed by Vim plugin
PATH=$(prepend-to-path "$HOME/.local/share/nvim/site/pack/packer/start/skim/bin")

export PATH
