# home/.config/profile.d/40-nix
if [ -e "$HOME/.nix-profile/etc/profile.d/nix.sh" ]; then
  source "$HOME/.nix-profile/etc/profile.d/nix.sh"
  export MANPATH="$HOME/.nix-profile/share/man:$MANPATH"
  export INFOPATH="$HOME/.nix-profile/share/info:$INFOPATH"
fi
