# home/.config/profile.d/40-nix
if [ -e "$HOME/.nix-profile/etc/profile.d/nix.sh" ]; then
  source "$HOME/.nix-profile/etc/profile.d/nix.sh"
fi

if [ -n "$NIX_LINK" ]; then
  export MANPATH="$NIX_LINK/share/man:$MANPATH"
  export INFOPATH="$NIX_LINK/share/info:$INFOPATH"
  export CMAKE_PREFIX_PATH="$NIX_LINK"
fi
