# home/.config/profile.d/60-docker-support
alias docker='sudo docker'
alias docker-gc='sudo ''docker ps -aq --filter "status=exited" | xargs -r docker rm && docker images -q --filter "dangling=true" | xargs -r docker rmi'''
