# home/.config/profile.d/60-elm
export PATH=$PATH:/usr/local/Elm-Platform/current/.cabal-sandbox/bin
