# home/.config/profile.d/60-ghcup
# ghcup manages installations of ghc, cabal, stack, haskell-language-server, etc

[ -f "/home/jesse/.ghcup/env" ] && source "/home/jesse/.ghcup/env" # ghcup-env
