# home/.config/profile.d/60-ghcup
# ghcup manages installations of ghc, cabal, stack, haskell-language-server, etc
# Check whether it's already added its path entry to avoid setting it up twice.
if [[ -f "/home/jesse/.ghcup/env" && ":$PATH:" != *":$HOME/.ghcup/bin:"* ]]; then
  source "/home/jesse/.ghcup/env"
fi
