# home/.config/profile.d/80-deno
export DENO_INSTALL="$HOME/.deno"
export PATH=$(prepend-to-path "$DENO_INSTALL/bin:$PATH")
