# home/.config/profile.d/80-fzf
export FZF_DEFAULT_COMMAND='rg --files --follow 2>/dev/null'
export FZF_CTRL_T_COMMAND='rg --files --follow 2>/dev/null'
