# home/.config/profile.d/80-postgres
if [ -d "/usr/lib/postgresql/11/bin" ]; then
  PATH=$(prepend-to-path "/usr/lib/postgresql/11/bin")
fi
