# home/.config/profile.d/80-rust
if [ -d "$HOME/.cargo/bin" ]; then
  PATH=$(prepend-to-path "$HOME/.cargo/bin")
fi
